# longan/test/dragonboard/rootfs/var/lib/dpkg/info/isdnutils-base.list
/.
/etc
/etc/init.d
/etc/init.d/isdnutils-base
/etc/isdn
/etc/isdn/init.d.functions
/usr
/usr/bin
/usr/bin/isdncause
/usr/lib
/usr/lib/libisdnudp.a
/usr/sbin
/usr/sbin/divertctrl
/usr/sbin/hisaxctrl
/usr/sbin/imon
/usr/sbin/imontty
/usr/sbin/iprofd
/usr/sbin/isdnconfig
/usr/sbin/isdnctrl
/usr/share
/usr/share/doc
/usr/share/doc/isdnutils-base
/usr/share/doc/isdnutils-base/HOWTO.gz
/usr/share/doc/isdnutils-base/NEWS.gz
/usr/share/doc/isdnutils-base/README.Debian
/usr/share/doc/isdnutils-base/README.HiSax.gz
/usr/share/doc/isdnutils-base/README.gz
/usr/share/doc/isdnutils-base/README.libipt
/usr/share/doc/isdnutils-base/changelog.Debian.gz
/usr/share/doc/isdnutils-base/copyright
/usr/share/doc/isdnutils-base/examples
/usr/share/doc/isdnutils-base/examples/event.d-ttyI0
/usr/share/doc/isdnutils-base/examples/grepconfig.pl
/usr/share/isdn
/usr/share/isdn/default
/usr/share/isdn/default/device.DEVICE
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/isdnutils-base
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/imon.1.gz
/usr/share/man/man1/imontty.1.gz
/usr/share/man/man4
/usr/share/man/man4/isdn_audio.4.gz
/usr/share/man/man4/isdnctrl.4.gz
/usr/share/man/man4/isdninfo.4.gz
/usr/share/man/man4/ttyI.4.gz
/usr/share/man/man7
/usr/share/man/man7/isdn_cause.7.gz
/usr/share/man/man8
/usr/share/man/man8/.isdnctrl_conf.8.gz
/usr/share/man/man8/divertctrl.8.gz
/usr/share/man/man8/hisaxctrl.8.gz
/usr/share/man/man8/iprofd.8.gz
/usr/share/man/man8/isdnconfig.8.gz
/usr/share/man/man8/isdnctrl.8.gz
/usr/share/menu
/usr/share/menu/isdnutils-base
/usr/share/doc/isdnutils-base/CONFIG.gz
/usr/share/doc/isdnutils-base/HOWTO.isdnutils.gz
/usr/share/man/man1/isdncause.1.gz
