# longan/test/dragonboard/rootfs/var/lib/dpkg/info/isdnutils-base.conffiles
/etc/init.d/isdnutils-base
/etc/isdn/init.d.functions
