# longan/test/dragonboard/rootfs/usr/lib/systemd/system/rcs.service
[Unit]
Description=rcS-Service
After=rcS.service
[Service]
Type=forking
ExecStart=/etc/init.d/rcS start
[Install]
WantedBy=multi-user.target

